# systemd/calf.service
[Service]
User=http
ExecStart=@pkglibexecdir@/calf.fcgi

[Install]
Also=calf.socket
