# systemd/calf.socket
[Socket]
ListenStream=/run/calf.sock

[Install]
WantedBy=sockets.target
