# systemd/calf.service.in
# You shouldn't start this unit directly, as calf will think it is run as CGI.
# Start calf.socket instead, or use something like spawn-fcgi.

[Service]
User=http
ExecStart=@pkglibexecdir@/calf.fcgi

[Install]
Also=calf.socket
